# systemd/cloud-sync-watchdog.timer
# SPDX-License-Identifier: MPL-2.0
[Unit]
Description=Cloud Sync Watchdog Timer - Periodic Health Check
Documentation=https://github.com/hyperpolymath/cloud-sync-tuner

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
